# backend/launchd/com.h5app.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.h5app.backend</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/ec2-user/h5-app/backend/scripts/start_backend.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ec2-user/h5-app/backend</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/h5-app-backend.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/h5-app-backend.log</string>
</dict>
</plist>
